# roles/web/templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream={{ app_socket }}

[Install]
WantedBy=sockets.target
